# myservice.service
[Service]
ExecStart=/usr/bin/usbDetectServ.sh

[Install]
WantedBy=multi-user.target
