# myservice.service
[Service]
ExecStart=bin/bash /usr/bin/usbDetectServ.sh

[Install]
WantedBy=multi-user.target
